# nanofused.service
[Unit]
Description=NanoFuse API Daemon
Documentation=https://github.com/jpoley/nanofuse
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nanofused
Restart=on-failure
RestartSec=5s

# Systemd will create /var/lib/nanofuse automatically
StateDirectory=nanofuse
StateDirectoryMode=0755

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security - minimal restrictions needed for Firecracker
NoNewPrivileges=false
ProtectHome=true
PrivateTmp=true

# Capabilities for networking
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN

# User/Group (runs as root for Firecracker and KVM access)
User=root
Group=root

[Install]
WantedBy=multi-user.target
