# nanofused.service
[Unit]
Description=NanoFuse API Daemon
Documentation=https://github.com/daax-dev/nanofuse
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nanofused
Restart=on-failure
RestartSec=5s

# Systemd will create /var/lib/nanofuse automatically
StateDirectory=nanofuse
StateDirectoryMode=0755

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security - minimal restrictions needed for Firecracker
NoNewPrivileges=false
ProtectHome=true
PrivateTmp=true

# Capabilities: networking (VM tap/NAT) + container-rootfs extraction.
# CAP_CHOWN/FOWNER/DAC_OVERRIDE/MKNOD/SETFCAP let the daemon reconstruct a pulled
# container's filesystem faithfully (e.g. chown etc/shadow to root:shadow, device
# nodes, file capabilities); without them tar extraction fails as EPERM even as root.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_MKNOD CAP_SETFCAP
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_MKNOD CAP_SETFCAP

# User/Group (runs as root for Firecracker and KVM access)
User=root
Group=root

[Install]
WantedBy=multi-user.target
